# multiContainer/from-github-multistage/apache/supervisord.conf
[supervisord]
nodaemon=true
user=root
loglevel=info
pidfile=/var/run/supervisord.pid ; (supervisord pidfile;default supervisord.pid)
logfile=/proc/self/fd/1
logfile_maxbytes = 0
;childlogdir=/var/log/ ; ('AUTO' child log dir, default $TEMP)

[unix_http_server]
file=/var/run/supervisor.sock   ; (the path to the socket file)
chmod=0700                       ; sockef file mode (default 0700)

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl = unix:///var/run/supervisor.sock ; use a unix:// URL  for a unix socket

[program:cron]
command = /usr/sbin/cron -f
user = root
autostart = true
startretries= 0
stdout_logfile = /var/log/supervisor/%(program_name)s.log
stderr_logfile = /var/log/supervisor/%(program_name)s.log

[program:apache2]
command = /usr/sbin/apache2ctl -DFOREGROUND
autorestart = false
killasgroup=true
stopasgroup=true
autostart=true

[program:atd]
command = /usr/sbin/atd
stdout_logfile = /var/log/supervisor/%(program_name)s.log
stderr_logfile = /var/log/supervisor/%(program_name)s.log
autostart = true
startretries = 1
start-secs=0

#[program:sshd]
#command=/usr/sbin/sshd -D
#autorestart=true